# cron.d/maho
# MAHO
*/5 *	* * *	root	/usr/share/maho/maho_monit.sh
